# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:laravel-server]
command=php artisan serve --host=0.0.0.0 --port=80
directory=/var/www/html
user=sail
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=3600

[program:laravel-echo-server]
command=bash -c "if [ -f laravel-echo-server.production.json ]; then laravel-echo-server start laravel-echo-server.production.json; elif [ -f laravel-echo-server.local.json ]; then laravel-echo-server start laravel-echo-server.local.json; else laravel-echo-server start; fi"
directory=/var/www/html
user=sail
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/laravel-echo-server.err.log
stdout_logfile=/var/log/supervisor/laravel-echo-server.out.log
stopwaitsecs=3600

# =============================================================================
# MQTT Consumer (PHP-MQTT nativo)
# Único punto de ingesta MQTT. Suscribe a todos los topics de dispositivos IoT:
# v1/mc/data, mc/data, v1/mc/alert, v1/mc/alert_control, v1/mc/ack, v1/mc/real_time
# =============================================================================
[program:mqtt-consumer]
command=php artisan mqtt:consume
directory=/var/www/html
user=sail
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/mqtt-consumer.err.log
stdout_logfile=/var/log/supervisor/mqtt-consumer.out.log
stopwaitsecs=3600

[program:queue-worker]
command=php artisan queue:work --queue=spot,spot1,spot4,default --sleep=3 --tries=3 --max-time=3600
directory=/var/www/html
user=sail
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/queue-worker.err.log
stdout_logfile=/var/log/supervisor/queue-worker.out.log
stopwaitsecs=3600
numprocs=1

[program:scheduler]
command=php artisan schedule:work
directory=/var/www/html
user=sail
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/scheduler.err.log
stdout_logfile=/var/log/supervisor/scheduler.out.log
stopwaitsecs=3600

# =============================================================================
# DESACTIVADOS: Scripts Python legacy (reemplazados por mqtt-consumer PHP)
# Los scripts Python eran el camino original (Mar 2022) para recibir datos MQTT.
# Fueron reemplazados por ConsumerCommand (PHP-MQTT) que es más eficiente
# (sin overhead HTTP) y cubre más topics (alerts, ack, config).
# Se mantienen comentados como referencia.
# =============================================================================
# [program:mqtt-receiver]
# command=python3 /usr/local/bin/receiveMqttEvent.py
# directory=/var/www/html
# user=sail
# autostart=true
# autorestart=true
# stderr_logfile=/var/log/supervisor/mqtt-receiver.err.log
# stdout_logfile=/var/log/supervisor/mqtt-receiver.out.log
# stopwaitsecs=3600

# [program:mqtt-realtime-receiver]
# command=python3 /usr/local/bin/receiveMqttRealTimeEvent.py
# directory=/var/www/html
# user=sail
# autostart=true
# autorestart=true
# stderr_logfile=/var/log/supervisor/mqtt-realtime-receiver.err.log
# stdout_logfile=/var/log/supervisor/mqtt-realtime-receiver.out.log
# stopwaitsecs=3600
